# systemd/nestcam-watchdog.service
[Unit]
Description=nestcam-streamer: frame-progress stall check (FR7/FR7b)
Documentation=file:///opt/PigeonCamSteward/SPEC.md

[Service]
Type=oneshot
User=root
Environment=NESTCAM_CONFIG=/etc/nestcam/config.yaml
ExecStart=/opt/PigeonCamSteward/bin/nestcam-watchdog.sh
SyslogIdentifier=nestcam-watchdog
