# systemd/nestcam-status-check.service
[Unit]
Description=nestcam-streamer: external YouTube live-status check (FR7c/d/e)
Documentation=file:///opt/nestcam-streamer/SPEC.md

[Service]
Type=oneshot
User=root
Environment=NESTCAM_CONFIG=/etc/nestcam/config.yaml
ExecStart=/opt/nestcam-streamer/bin/nestcam-status-check.sh
SyslogIdentifier=nestcam-status-check
