# systemd/pigeoncam-stream.service
[Unit]
Description=PigeonCamSteward: capture, encode, and push the live stream to YouTube
Documentation=file:///opt/PigeonCamSteward/SPEC.md
After=network-online.target
Wants=network-online.target

# FR6: load-bearing, not boilerplate. systemd's default start-limit (5
# starts within 10s) permanently stops restarting a unit that fails quickly
# in succession - exactly what happens when the camera is unplugged or
# wedged and ffmpeg exits immediately on launch. Disabling it here
# (combined with RestartSec below) means a 3am USB glitch recovers on its
# own instead of leaving the service dead for the rest of the season.
# pigeoncam-doctor.sh (FR17) checks that this line is present.
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Environment=PIGEONCAM_CONFIG=/etc/pigeoncam/config.yaml
ExecStart=/opt/PigeonCamSteward/bin/pigeoncam-stream.sh
Restart=always
RestartSec=10
SyslogIdentifier=pigeoncam-stream

# Hardening (C3). Deliberately compatible with two things this unit
# specifically needs that a blanket lockdown would break: full /dev
# access for the capture device (no PrivateDevices=/DeviceAllow= -
# camera.device is config-driven, not a fixed set this file could
# enumerate) and /run/user access for the audio.mode:real cross-user
# PipeWire/PulseAudio bridge (lib/pigeoncam-common.sh's
# resolve_pulse_bridge_env) - so no ProtectHome= here, unlike the other
# units. MemoryDenyWriteExecute= is deliberately not set anywhere in this
# project: ffmpeg builds have a real history of tripping it (SIMD
# dispatch/codec-path incompatibilities), and this is the one unit where
# a hardening-induced encode failure would be worst.
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
