# systemd/pigeoncam-archive-trim.timer
[Unit]
Description=Hourly trigger for pigeoncam-archive-trim.service (FR11)

[Timer]
# Fires shortly after each hour's segment(s) close - segments are cut on
# clock-hour boundaries (pigeoncam-stream.sh's segment_atclocktime=1).
OnCalendar=*-*-* *:05:00
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
