# systemd/pigeoncam-archive-trim.service
[Unit]
Description=PigeonCamSteward: hourly archive retention pass (FR11)
Documentation=file:///opt/PigeonCamSteward/SPEC.md

[Service]
Type=oneshot
User=root
Environment=PIGEONCAM_CONFIG=/etc/pigeoncam/config.yaml
ExecStart=/opt/PigeonCamSteward/bin/pigeoncam-archive-trim.sh
SyslogIdentifier=pigeoncam-archive-trim
TimeoutStartSec=300
